# systemd/git-backup.timer
[Unit]
Description=Daily Git commit+push for BookStack repo

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
